# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/deploy/bin/production.py
