# examples/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=/env/bin/python /deploy/bin/production.py
