# examples/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python /deploy/bin/production.py
